# tools/ubuntu/va-connect-watchdog-web.service
[Unit]
Description=VA-Connect watchdog web UI
After=network-online.target va-connect-site-watchdog.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
Environment=SITE_WATCHDOG_CONFIG=/opt/va-connect-watchdog/site-watchdog.json
ExecStart=/usr/bin/python3 /opt/va-connect-watchdog/va_connect_watchdog_web.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
