# tools/ubuntu/deploy/site_watchdog.service
[Unit]
Description=VA-Connect v2 site watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/va-connect-watchdog
Environment=PYTHONUNBUFFERED=1
Environment=VA_CONNECT_V2_CONFIG=/opt/va-connect-watchdog/site-watchdog.json
Environment=VA_CONNECT_V2_DATA_DIR=/var/lib/va-connect-site-watchdog
EnvironmentFile=-/opt/va-connect-watchdog/site_watchdog.env
ExecStart=/usr/bin/python3 -m tools.ubuntu.runtime.site_watchdog
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
